# contrib/init/systemd/mirrorbits.service.in
[Unit]
Description=Mirrorbits redirector
Documentation=https://github.com/etix/mirrorbits
After=network.target

[Service]
DynamicUser=yes
LogsDirectory=mirrorbits
RuntimeDirectory=mirrorbits
PIDFile=/run/mirrorbits/mirrorbits.pid
ExecStart=##PREFIX##/bin/mirrorbits -D -p /run/mirrorbits/mirrorbits.pid
ExecReload=/bin/kill -HUP $MAINPID
ExecStop=-/bin/kill -QUIT $MAINPID
TimeoutStopSec=5
KillMode=mixed
Restart=on-failure

[Install]
WantedBy=multi-user.target
