# contrib/init/systemd/mirrorbits.service
[Unit]
Description=Mirrorbits redirector

[Service]
ExecStart=/usr/bin/mirrorbits -D